# citizenpassport/deploy/linux/systemd/citizenpassport-backup.timer
[Unit]
Description=Run CitizenPassport backup daily

[Timer]
OnCalendar=*-*-* 02:15:00
Persistent=true
RandomizedDelaySec=300
Unit=citizenpassport-backup.service

[Install]
WantedBy=timers.target
