# citizenchain/scripts/citizenchain-node.service
[Unit]
Description=CitizenChain Node
After=network.target

[Service]
Type=simple
User=citizenchain
ExecStart=/opt/citizenchain/node \
  --chain dev \
  --base-path /opt/citizenchain/data \
  --port 30333 \
  --rpc-port 9944 \
  --rpc-methods Unsafe \
  --rpc-cors all \
  --no-prometheus \
  --listen-addr /ip4/0.0.0.0/tcp/30333/ws
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
