# CPMS/deploy/linux/systemd/cpms-backend.service
[Unit]
Description=CPMS Backend Service
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
User=cpms
Group=cpms
WorkingDirectory=/var/lib/cpms
EnvironmentFile=/etc/cpms/cpms-backend.env
ExecStart=/opt/cpms/bin/cpms-backend
Restart=always
RestartSec=3
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
